# firmware_analysis/Firmware/_D-link.bin.extracted/squashfs-root/usr/lib/opkg/info/firewall.list
/etc/firewall.acl
/etc/firewall.guest
/etc/firewall.adv
/etc/init.d/firewall
/sbin/fw3
/etc/config/acl
/etc/firewall.parentcontrol
/etc/config/firewall
/etc/init.d/fix_firewall
/lib/upgrade/keep.d/firewall
/etc/hotplug.d/iface/20-firewall
/etc/config/parentcontrol
/etc/firewall.portal
/etc/config/AdvFirewall
