# firmware_analysis/Firmware/_D-link.bin.extracted/squashfs-root/usr/lib/opkg/info/reptd.list
/etc/config/reptd
/usr/sbin/reptd
/etc/init.d/reptd
